# systemd/rezept-daily-import.service
[Unit]
Description=Rezept Tagebuch Daily TheMealDB Import (Vegetarian)
After=network.target seaser-rezept-tagebuch.service

[Service]
Type=oneshot
EnvironmentFile=-/home/gabor/easer_projekte/rezept-tagebuch/.env
ExecStart=/home/gabor/easer_projekte/rezept-tagebuch/scripts/daily-import.sh by_category Vegetarian
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
