# systemd/quadlets/dev/seaser-postgres-dev.container
# PostgreSQL Database for Rezept-Tagebuch DEV
# Network: seaser-network
# Database: rezepte_dev

[Unit]
Description=PostgreSQL Database for Rezept-Tagebuch DEV
Documentation=https://github.com/easer/rezept-tagebuch

[Container]
Image=docker.io/library/postgres:16-alpine
ContainerName=seaser-postgres-dev
Network=seaser-network

# Database Configuration
Environment=POSTGRES_PASSWORD=seaser
Environment=POSTGRES_DB=rezepte_dev
Environment=POSTGRES_USER=postgres

# Data Persistence
Volume=/home/gabor/easer_projekte/rezept-tagebuch/data/postgres-dev:/var/lib/postgresql/data:Z

[Service]
Restart=on-failure
TimeoutStopSec=70

[Install]
WantedBy=default.target
