# systemd/quadlets/dev/seaser-rezept-tagebuch-dev.container
# Rezept-Tagebuch DEV App Container
# Network: seaser-network
# URL: http://192.168.2.139:8000/rezept-tagebuch-dev/

[Unit]
Description=Rezept-Tagebuch DEV App
Documentation=https://github.com/easer/rezept-tagebuch
After=seaser-postgres-dev.service
Wants=seaser-postgres-dev.service

[Container]
Image=localhost/seaser-rezept-tagebuch:dev
ContainerName=seaser-rezept-tagebuch-dev
Network=seaser-network

# Application Configuration
Environment=APP_VERSION=dev
Environment=DB_TYPE=postgresql
Environment=DEV_MODE=true
Environment=POSTGRES_HOST=seaser-postgres-dev
Environment=POSTGRES_PASSWORD=seaser
Environment=POSTGRES_DB=rezepte_dev
Environment=POSTGRES_USER=postgres

# Load additional environment from file (DeepL API Key, etc.)
EnvironmentFile=/home/gabor/easer_projekte/rezept-tagebuch/.env

# Data Persistence
Volume=/home/gabor/easer_projekte/rezept-tagebuch/data/dev:/data:Z

[Service]
Restart=on-failure
TimeoutStopSec=70

[Install]
WantedBy=default.target
